# cookbooks/zookeeper/templates/zookeeper.service.erb
[Unit]
Description=Zookeeper Service
After=network.target
 
[Service]
User=root
Group=root
Type=simple
ExecStart=/opt/zookeeper/apache-zookeeper-3.5.6-bin/bin/zkServer.sh start-foreground
 
[Install]
WantedBy=multi-user.target
